# deploy/xszn-comments-backup.timer
[Unit]
Description=Daily backup for the NJUPT freshman guide comments database

[Timer]
OnCalendar=*-*-* 04:20:00
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
